# deploy/lastshot-tracker.service
[Unit]
Description=last_shot tracker + shadow executor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/last-shot
ExecStart=/home/ubuntu/last-shot/.venv/bin/python -m src.tracking.run --status-every 300
Restart=always
RestartSec=10
# tracker logs are chatty; cap journal usage
LogRateLimitIntervalSec=30
LogRateLimitBurst=2000

[Install]
WantedBy=multi-user.target
